# scripts/systemd/ai-nodes-health.service
[Unit]
Description=AI Nodes Health Check Service
Documentation=https://github.com/your-org/ai-nodes
After=network.target
Wants=network.target

[Service]
Type=forking
User=ai-nodes
Group=ai-nodes
WorkingDirectory=/opt/ai-nodes
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
EnvironmentFile=/opt/ai-nodes/.env

# Health check every 5 minutes
ExecStart=/bin/bash -c 'while true; do /opt/ai-nodes/scripts/health-checks/basic.sh; sleep 300; done'

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ai-nodes/logs /opt/ai-nodes/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-nodes-health

# Process management
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target